# config/hyprwhspr.service
[Unit]
Description=hyprwhspr stt
Documentation=https://github.com/goodroot/hyprwhspr
After=graphical-session.target pipewire.service ydotool.service
StartLimitIntervalSec=300
StartLimitBurst=5
Requires=ydotool.service
Wants=pipewire.service

[Service]
Type=simple
# Wait for Wayland compositor - check for socket file directly since systemd
# services don't inherit WAYLAND_DISPLAY from the graphical session
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do ls /run/user/$(id -u)/wayland-* >/dev/null 2>&1 && exit 0; sleep 0.5; done; echo "Warning: Wayland socket not found after 15s"'
ExecStart=/usr/lib/hyprwhspr/bin/hyprwhspr
# Ensure virtual keyboard is killed even on crash/SIGKILL to release device grabs
ExecStopPost=/bin/bash -c 'pkill -9 -f "hyprwhspr-virtual-keyboard" 2>/dev/null || true'
Environment=HYPRWHSPR_ROOT=/usr/lib/hyprwhspr
Environment=PYTHONUNBUFFERED=1
Environment=LD_LIBRARY_PATH=/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cublas/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cudnn/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cufft/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cuda_runtime/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cuda_nvrtc/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cusparse/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cusolver/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/curand/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/nvjitlink/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/nccl/lib:/opt/cuda/lib64
Restart=on-failure
# Give kernel time to release device grabs after stop
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target