# config/hyprwhspr.service
[Unit]
Description=hyprwhspr stt
Documentation=https://github.com/goodroot/hyprwhspr
After=graphical-session.target pipewire.service ydotool.service
StartLimitIntervalSec=300
StartLimitBurst=5
Requires=ydotool.service
Wants=pipewire.service

[Service]
Type=simple
# Pass Wayland/X11 display env from the user session manager so wl-copy,
# wl-paste, and the GTK4 mic-osd overlay can connect to the compositor.
PassEnvironment=WAYLAND_DISPLAY DISPLAY
# Wait for Wayland compositor socket to appear before starting
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do ls /run/user/$(id -u)/wayland-* >/dev/null 2>&1 && exit 0; sleep 0.5; done; echo "Warning: Wayland socket not found after 15s"'
# Wait for UWSM to export WAYLAND_DISPLAY into the systemd user environment.
# The socket may exist before UWSM calls import-environment, so PassEnvironment
# would pass an empty value without this guard — leaving subprocesses (mic-osd,
# wl-copy) unable to connect to the compositor.
ExecStartPre=/bin/bash -c 'for i in $(seq 1 60); do systemctl --user show-environment 2>/dev/null | grep -q "^WAYLAND_DISPLAY=" && exit 0; sleep 0.5; done; echo "Warning: WAYLAND_DISPLAY not found in session env after 30s"'
# Clean stale mic-osd PID file so the runner spawns a fresh daemon
# instead of trying to signal a PID that was reused by another process
ExecStartPre=/bin/rm -f %h/.config/hyprwhspr/mic_osd.pid
ExecStart=/usr/lib/hyprwhspr/bin/hyprwhspr
# Ensure virtual keyboard is killed even on crash/SIGKILL to release device grabs
ExecStopPost=/bin/bash -c 'pkill -9 -f "hyprwhspr-virtual-keyboard" 2>/dev/null || true'
Environment=HYPRWHSPR_ROOT=/usr/lib/hyprwhspr
Environment=PYTHONUNBUFFERED=1
Environment=LD_LIBRARY_PATH=/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cublas/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cudnn/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cufft/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cuda_runtime/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cuda_nvrtc/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cusparse/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/cusolver/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/curand/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/nvjitlink/lib:/home/buddleman/.local/share/hyprwhspr/venv/lib/python3.14/site-packages/nvidia/nccl/lib:/opt/cuda/lib64
Restart=on-failure
# Give kernel time to release device grabs after stop
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target